# packages/zymbit/zymbit-0.5.17.tar.gz/zymbit-0.5.17/files/raspberrypi-arch/etc/systemd/system/zymbit.service
[Install]
WantedBy=multi-user.target

[Unit]
Description=Zymbit Cloud Connector
After=network.target

[Service]
Restart=always
RestartSec=5
ExecStart=/opt/zymbit/latest/bin/zymbit connect
