# packages/FERG/FERG-0.9.5.tar.gz/FERG-0.9.5/comptage.service
[Unit]
Description= FERG : Mon vieux compteur EDF connecte en mode framboise
After=network-online.target

[Service]
ExecStart=/opt/FERG/comptage.sh
User=root
KillMode=process

[Install]
WantedBy=multi-user.target
